# supervisor/queue_dashboard.ini
[program:queue_dashboard]                                                                                                                                                                                
directory=opt/asterisk_queue_dashboard
command=python run.py
user=root 
autorestart=true
