# supervisor/queue_dashboard.ini
[program:queue_dashboard]                                                                                                                                                                                
directory=/home/andriy
command=/home/andriy/queue_log.sh
user=root 
autorestart=true
 